# systemd_2/user/hyprland-session.target
[Unit]
Description=Hyprland Graphical Session
BindsTo=graphical-session.target

